# deploy/bytebit-backend.service
[Unit]
Description=ByteBit Backend (Daphne ASGI)
After=network.target redis-server.service

[Service]
Type=simple
User=bytebit
Group=bytebit
WorkingDirectory=/opt/bytebit/app/ByteBit-backend

# Load all env vars from the secrets file
EnvironmentFile=/opt/bytebit/.env

ExecStart=/opt/bytebit/venv/bin/daphne \
    -b 127.0.0.1 \
    -p 8000 \
    --access-log - \
    backend.asgi:application

Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bytebit-backend

[Install]
WantedBy=multi-user.target
